# deploy/systemd/hfi-restore-check.service
[Unit]
Description=HFI monthly backup restore validation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/hfi/app
Environment=ENV_FILE=/opt/hfi/app/.env.prod
ExecStart=/usr/bin/env bash /opt/hfi/app/deploy/scripts/restore_check.sh
