# deploy/systemd/hfi-restore-check.timer
[Unit]
Description=Run HFI monthly restore validation

[Timer]
OnCalendar=monthly
Persistent=true
Unit=hfi-restore-check.service

[Install]
WantedBy=timers.target
