# deploy/systemd/tailscale-serve.service
[Unit]
Description=Expose HFI privately over Tailscale Serve
Requires=tailscaled.service hfi-prod.service
After=tailscaled.service hfi-prod.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/tailscale serve --bg --https=443 http://127.0.0.1:80
ExecStop=/usr/bin/tailscale serve --https=443 off

[Install]
WantedBy=multi-user.target
